# docker/etc/services.d/horizon/run
#!/command/execlineb -P

/var/www/foxes-knowledge/artisan horizon
